# deployments/hermes-bitnet-lxc/systemd/bitnet-llama-server.service.example
# User systemd unit — install to ~/.config/systemd/user/bitnet-llama-server.service
# Set BITNET_GGUF to your .gguf path in an EnvironmentFile.
#
#   mkdir -p ~/.config/environment.d
#   echo 'BITNET_GGUF=/home/you/src/BitNet/models/.../ggml-model-i2_s.gguf' >> ~/.config/environment.d/bitnet.conf
#
# Or use drop-in: systemctl --user edit bitnet-llama-server

[Unit]
Description=BitNet llama-server (OpenAI-compatible API for Hermes)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/src/BitNet
Environment=BITNET_DIR=%h/src/BitNet
Environment=BITNET_BUILD_DIR=%h/src/BitNet/build
Environment=BITNET_PORT=8080
# EnvironmentFile=-%h/.config/hermes-bitnet-lxc/bitnet.env
# After bootstrap, run-bitnet-server is symlinked to ~/.local/bin (see bootstrap-all.sh).
ExecStart=%h/.local/bin/run-bitnet-server.sh
Restart=on-failure
RestartSec=10
# Large models: raise or remove after tuning for your LXC RAM
MemoryMax=48G

[Install]
WantedBy=default.target
