# deployments/hermes-bitnet-lxc/systemd/code-server.service.example
# User systemd unit — install to ~/.config/systemd/user/code-server.service
# then: systemctl --user daemon-reload && systemctl --user enable --now code-server
#
# Adjust paths and ports. Prefer SSH -L over exposing 0.0.0.0.
# Default qminiwasm-core clone: C:\GitHub\LLM_Pract\qminiwasm-core (WSL: /mnt/c/GitHub/LLM_Pract/qminiwasm-core).

[Unit]
Description=code-server (VS Code web) — qminiwasm-core primary workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HERMES_BITNET_SRC_ROOT=%h/src
# Use your qminiwasm tree: WSL path below matches C:\GitHub\LLM_Pract\qminiwasm-core. On Incus without /mnt/c, use e.g. %h/src/qminiwasm-core.
ExecStart=%h/.local/bin/code-server --bind-addr 127.0.0.1:8443 /mnt/c/GitHub/LLM_Pract/qminiwasm-core
Restart=on-failure
RestartSec=5
MemoryMax=2G

[Install]
WantedBy=default.target
