# deployments/hermes-bitnet-lxc/systemd/openvscode-server.service.example
# User systemd unit — install to ~/.config/systemd/user/openvscode-server.service
#
# Default (this file): no connection token — open http://127.0.0.1:3010/ directly.
# Binds OPENVS_CODE_BIND (127.0.0.1) so the editor is not offered tokenless on all interfaces.
# If your browser cannot reach WSL on 127.0.0.1, use: systemctl --user edit openvscode-server
#   and set ExecStart to use --host 0.0.0.0 (still without token = only for trusted LAN).
#
# Token mode (LAN / shared host): use openvscode-server-token.service.example instead, or set
# OPENVS_CODE_REQUIRE_TOKEN=1 when running 07-openvscode-server.sh.
#
#   mkdir -p ~/.config/systemd/user
#   cp systemd/openvscode-server.service.example ~/.config/systemd/user/openvscode-server.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now openvscode-server
#   loginctl enable-linger "$USER"

[Unit]
Description=OpenVSCode Server (browser IDE, localhost, no URL token)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HOME=%h
# So the integrated terminal and ACP-launched tools find hermes, uv, etc.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
# BitNet llama-server rejects OpenAI ``tools``; Hermes honors this env (patched run_agent.py via 08-hermes-bitnet-config.sh).
Environment=HERMES_CHAT_COMPLETIONS_NO_TOOLS=1
WorkingDirectory=%h
ExecStart=%h/.local/bin/openvscode-server --host 127.0.0.1 --port 3010 --without-connection-token --accept-server-license-terms
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
